# root/etc/s6-overlay/s6-rc.d/novnc/run
#!/command/with-contenv bash

exec \
    s6-setuidgid twgs \
        /usr/bin/novnc_server --vnc localhost:5900 --listen 8080