# root/etc/s6-overlay/s6-rc.d/x11vnc/run
#!/command/with-contenv bash

exec s6-setuidgid abc \
        /usr/bin/x11vnc -noxrecord