# root/etc/s6-overlay/s6-rc.d/xvfb/run
#!/command/with-contenv bash

exec \
    s6-setuidgid abc \
        /usr/bin/Xvfb :0 -screen 0 1024x768x24 -nolisten unix