# root/etc/s6-overlay/s6-rc.d/windowmanager/run
#!/command/with-contenv bash

export HOME=/app

s6-setuidgid twgs /usr/bin/icewm-session