# root/etc/s6-overlay/s6-rc.d/windowmanager/run
#!/command/execlineb -P

with-contenv

unexport HOME

s6-setuidgid twgs /usr/bin/icewm-session