# root/etc/s6-overlay/s6-rc.d/twgs/run
#!/command/with-contenv bash

exec \
    s6-setuidgid twgs \
        /usr/bin/wine /app/drive_c/EIS/TWGS/twgs.exe